# jira-report.service
[Unit]
Description=Jira Report Web Application
After=network.target

[Service]
Type=simple
User=kit02
Group=kit02
WorkingDirectory=/home/kit02/jira_report
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/kit02/jira_report/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jira-report

[Install]
WantedBy=multi-user.target
